# debian/linux-config-manager.service
[Unit]
Description=Linux Config Manager Backend Service
After=network.target

[Service]
Type=simple
User=config-manager
Group=config-manager
WorkingDirectory=/opt/linux-config-manager
ExecStart=/opt/linux-config-manager/bin/config-manager-backend
Restart=always
RestartSec=5
Environment=PORT=8080
Environment=STATIC_DIR=/opt/linux-config-manager/web

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home /etc /var
PrivateTmp=true

[Install]
WantedBy=multi-user.target